# etc/supervisor/conf.d/sqs-init.conf
[program:sqs-init]
command=/opt/sqs-init.sh
priority=1
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

